# gluon-txpowerfix/files/etc/init.d/txpowerfix/txpowerfix
#!/bin/sh /etc/rc.common

START=99

boot() {
/lib/gluon/txpowerfix/txpowerfix.lua
rm /etc/init.d/txpowerfix
}